# config/supervisor/cron.conf
[program:cron]
priority=20
directory=/tmp
command=/usr/sbin/cron -f
autorestart=true