# config/supervisor/postfix.conf
[program:postfix]
process_name=master
priority=20
command=/usr/sbin/postfix -c /etc/postfix start
startsecs=0
autorestart=false