# assets/config/supervisor/redis.conf
[program:redis]
priority=10
command=/usr/bin/redis-server /etc/redis/redis.conf --daemonize no
user=redis
autorestart=true
