# assets/config/supervisor/rspamd.conf
[program:rspamd]
priority=10
command=/usr/bin/rspamd -f -c /etc/rspamd/rspamd.conf -u _rspamd -g _rspamd
user=_rspamd
autorestart=true
