# assets/config/supervisor/rsyslog.conf
[program:rsyslog]
priority=5
command=/usr/sbin/rsyslogd -n
autorestart=true