# assets/config/supervisor/spamass-milter.conf
[program:spamass-milter]
process_name=spamass-milter
priority=10
command=/etc/init.d/spamass-milter start
startsecs=0
autorestart=false