# assets/config/supervisor/spamassassin.conf
[program:spamassassin]
priority=10
command=/usr/sbin/spamd --pidfile=/var/run/spamd.pid --create-prefs --max-children 5 --helper-home-dir -x -u vmail
autorestart=true